# com.klublepszezycie.email.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.klublepszezycie.email</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Volumes/Dane/Projekty/devs/klublepszezycie/process_emails.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Volumes/Dane/Projekty/devs/klublepszezycie</string>
    <key>StartInterval</key>
    <integer>60</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/email_processor.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/email_processor_error.log</string>
</dict>
</plist>



